# deploy/systemd/send_notifications.timer
[Unit]
Description=Hourly timer to run send_notifications for Hendoshi

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
